# deploy/systemd/dquery-rules-update.service
[Unit]
Description=Update dquery ChinaMax rules
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=HOME=/root
Environment=GH_CONFIG_DIR=/root/.config/gh
ExecStart=/opt/dquery/deploy/update-rules.sh
